# packages/third/tools/acelist/system.d/acelist.service
[Unit]
Description=AceList Script
Requires=network-online.target
After=network-online.target
Before=kodi.service

[Service]
ExecStart=/usr/bin/python3 /usr/bin/acelist.py
ExecStop=/bin/kill -TERM $MAINPID

[Install]
WantedBy=kodi.target
